# sys-apps/ps3vram-swap/files/ps3vram-swap-0.0.1
#!/sbin/openrc-run

start() {
    ebegin "Starting ps3vram-swap"
    if ! swapon -s | grep -q '/dev/ps3vram'; then
        wipefs -a /dev/ps3vram > /dev/null
        mkswap /dev/ps3vram > /dev/null
        swapon /dev/ps3vram -p10 > /dev/null
    fi
    eend $?
}

stop() {
    ebegin "Stopping ps3vram-swap"
    if swapon -s | grep -q '/dev/ps3vram'; then
        swapoff /dev/ps3vram > /dev/null
    fi
    eend $?
}
